# crontab
0 0 * * * cd /code/scripts && python caribbean_jobs_scraper.py
